# templates/centos/init.erb
### BEGIN INIT INFO
# Provides:
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start daemon at boot time
# Description:       Enable service provided by daemon.
### END INIT INFO
dir="<%= @sbin_dir %>"
user="pingfederate"
name=`basename $0`

case "$1" in
    start)
        <% if %x(rpm -qa | grep systemd) != '' %>
        echo "Redirecting to systemd"
        /usr/bin/systemctl start <%= @service_name %>
        <% else %>
        su - $user -c $dir/pingfederate-run.sh
        <% end %>
    ;;
    stop)
        <% if %x(rpm -qa | grep systemd) != '' %>
        echo "Redirecting to systemd"
        /usr/bin/systemctl stop <%= @service_name %>
        <% else %>
        su - $user -c $dir/pingfederate-shutdown.sh
        <% end %>
    ;;
    restart)
        <% if %x(rpm -qa | grep systemd) != '' %>
        echo "Redirecting to systemd"
        /usr/bin/systemctl restart <%= @service_name %>
        <% else %>
	    $0 stop
	    sleep 3
	    $0 start
	   <% end %>
    ;;
    status)
        <% if %x(rpm -qa | grep systemd) != '' %>
        /usr/bin/systemctl status <%= @service_name %>
        echo "Redirecting to systemd"
        <% else %>
	    echo "Status is not available on the OS"
	   <% end %>
    ;;
    *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 1
    ;;
esac

exit 0